# smtp-probe-service/smtp-probe.service
[Unit]
Description=Ellyn SMTP Probe Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/smtp-probe-service
ExecStart=/usr/bin/node index.js
Restart=on-failure
RestartSec=5
EnvironmentFile=/home/ubuntu/smtp-probe-service/.env
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
